# deploy/sbb-kiosk-ordering.service
[Unit]
Description=SBB Kiosk Ordering
After=network.target

[Service]
Type=simple
WorkingDirectory=/root/SBB-Kiosk-Ordering
Environment=NODE_ENV=production
Environment=PORT=4110
Environment=PREVIEW_PORT=5174
EnvironmentFile=-/root/SBB-Kiosk-Ordering/.env
ExecStart=/usr/bin/npm run serve:all
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
